# 9a8d1991-2ecf-4e70-b77a-f506e811880b.myprj/modules/fastapi-tutorial/ansible/roles/fastapi_tutorial/templates/fastapi-tutorial.service.j2
[Unit]
Description=FastAPI Tutorial Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fastapi-tutorial
Environment="PATH=/opt/fastapi-tutorial/venv/bin"
ExecStart=/opt/fastapi-tutorial/venv/bin/uvicorn app.main:app --host {{ app_host }} --port {{ app_port }}
Restart=always

[Install]
WantedBy=multi-user.target