# 9a8d1991-2ecf-4e70-b77a-f506e811880b.myprj/modules/fastapi-tutorial/ansible/roles/fastapi_tutorial/templates/fastapi-tutorial.service.j2
[Unit]
Description=FastAPI Tutorial Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory={{ fastapi_app_dir }}
Environment="PATH={{ fastapi_python_venv }}/bin"
ExecStart={{ fastapi_python_venv }}/bin/uvicorn app.main:app --host {{ fastapi_app_host }} --port {{ fastapi_app_port }}
Restart=always

[Install]
WantedBy=multi-user.target