# 9a8d1991-2ecf-4e70-b77a-f506e811880b.myprj/modules/fastapi-tutorial/ansible/roles/fastapi_tutorial/templates/fastapi-tutorial.service.j2
[Unit]
Description=FastAPI Tutorial Service
After=network.target postgresql.service

[Service]
Type=simple
User={{ fastapi_service_user }}
WorkingDirectory={{ fastapi_app_path }}
Environment="PATH={{ fastapi_app_path }}/venv/bin"
ExecStart={{ fastapi_app_path }}/venv/bin/uvicorn app.main:app --host {{ fastapi_host }} --port {{ fastapi_port }}
Restart=always

[Install]
WantedBy=multi-user.target