# 38079264-0f1d-4866-a75e-52a53510d09b.myprj/modules/fastapi-tutorial/ansible/roles/fastapi_tutorial/templates/fastapi-tutorial.service.j2
[Unit]
Description=FastAPI Tutorial Service
After=network.target postgresql.service

[Service]
Type=simple
User={{ fastapi_service_user }}
WorkingDirectory={{ fastapi_app_dir }}
Environment="PATH={{ fastapi_venv_path }}/bin"
ExecStart={{ fastapi_venv_path }}/bin/uvicorn app.main:app --host 0.0.0.0 --port {{ fastapi_port }}
Restart=always

[Install]
WantedBy=multi-user.target